# configs-tunnel/supervisor/conf.d/rtunnel.conf
[program:rtunnel]
command=ssh -R 2222:localhost:2222 ubuntu@michaeljmiller.net
user=commnerd
