# configs/supervisor/conf.d/controller.conf
[program:controller]
command=/usr/sbin/workspace-controller
user=root
