# configs/supervisor/conf.d/dockerd.conf
[program:dockerd]
command=dockerd --experimental
user=root
