# configs/supervisor/conf.d/dockerd.conf
[program:dockerd]
command=dockerd --experimental -s overlay
user=root
