# configs/supervisor/conf.d/dockerd.conf
[program:dockerd]
command=dockerd
user=root
