# configs/supervisor/conf.d/nginx.conf
[program:nginx]
command=nginx -g 'daemon off;'
user=root
