# configs/supervisor/conf.d/rtunnel.conf
[program:rtunnel]
command=ssh -o ServerAliveInterval=300 -NR 2222:localhost:22 ubuntu@michaeljmiller.net || sleep 900
user=commnerd
autostart=false
startretries=32
